# systemd/ddns.service
[Unit]
Description=Dynamically update DNS records.
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/ddns.sh /etc/ddns/*
