# systemd/ddns.service
[Unit]
Description=Dynamically update DNS records.

[Service]
ExecStart=/usr/bin/ddns run
