# systemd/ddns.timer
[Unit]
Description=Timer for updating dynamic DNS records.

[Timer]
OnCalendar=*-*-* *:00/5:00
