# systemd/ddns.service
[Unit]
Description=Dynamically update DNS records.
Wants=ddns.timer

[Service]
ExecStart=/usr/bin/ddns run

[Install]
WantedBy=multi-user.target
