# systemd/ddns.timer
[Unit]
Description=Timer for updating dynamic DNS records.
PartOf=ddns.service

[Timer]
OnUnitActiveSec=5min
