# systemd/ddns.timer
[Unit]
Description=Timer for updating dynamic DNS records.

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=multi-user.target
